# packaging/dashboard/deb/postinst
#!/bin/sh
# Debian/Ubuntu post-install for the dashboard: create the service user, seed env.
set -e

case "$1" in
    configure)
        if ! getent passwd easy-failover-dashboard >/dev/null 2>&1; then
            useradd --system --no-create-home --shell /usr/sbin/nologin \
                --comment "easy-failover dashboard" easy-failover-dashboard || true
        fi
        if [ ! -e /etc/easy-failover-dashboard/dashboard.env ] \
            && [ -e /etc/easy-failover-dashboard/dashboard.env.example ]; then
            cp /etc/easy-failover-dashboard/dashboard.env.example \
                /etc/easy-failover-dashboard/dashboard.env
            chmod 0644 /etc/easy-failover-dashboard/dashboard.env
        fi
        if command -v systemctl >/dev/null 2>&1; then
            systemctl daemon-reload || true
        fi
        # Not enabled/started automatically: review dashboard.env (listen address,
        # roster, onboarding gate) then `systemctl enable --now easy-failover-dashboard`.
        ;;
esac

exit 0
