# packaging/systemd/easy-failover.service.in
[Unit]
Description=easy-failover virtual IP failover agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=@EASY_FAILOVER_SYSTEMD_BINDIR@/easy-failover --config @EASY_FAILOVER_SYSTEMD_SYSCONFDIR@/easy-failover/config.toml
Restart=always
RestartSec=2s

# TODO: Add least-privilege hardening once real VIP movement exists.
# TODO: See docs/linux-capabilities.md before changing granted capabilities.
# TODO: Use CAP_NET_ADMIN and CAP_NET_RAW instead of broad root privileges where possible.
# TODO: Evaluate ProtectSystem=strict with explicit writable paths.
# TODO: Enable NoNewPrivileges=true after validating required capabilities.

[Install]
WantedBy=multi-user.target
