# packaging/systemd/easy-failover.service.in
[Unit]
Description=easy-failover virtual IP failover agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=@EASY_FAILOVER_SYSTEMD_BINDIR@/easy-failover --config @EASY_FAILOVER_SYSTEMD_SYSCONFDIR@/easy-failover/config.toml --run-forever
Restart=on-failure
RestartSec=2s

# Baseline hardening keeps host networking available for VIP and heartbeat operations.
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
NoNewPrivileges=true
ProtectSystem=strict
# The write API rewrites config.toml (+ a .bak) in place on apply, so the config dir must
# stay writable under ProtectSystem=strict. Everything else remains read-only.
ReadWritePaths=@EASY_FAILOVER_SYSTEMD_SYSCONFDIR@/easy-failover
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK AF_PACKET
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
