# packaging/deb/postinst
#!/bin/sh
# Debian/Ubuntu post-install: seed config, refresh systemd, set up the bundled dashboard.
set -e

case "$1" in
    configure)
        # On first install, create an active config from the shipped example if the
        # operator does not already have one. Upgrades never overwrite it.
        if [ ! -e /etc/easy-failover/config.toml ] \
            && [ -e /etc/easy-failover/config.example.toml ]; then
            cp /etc/easy-failover/config.example.toml /etc/easy-failover/config.toml
            chmod 0644 /etc/easy-failover/config.toml
        fi

        # Bundled dashboard (present only when the package was built with it): create its
        # unprivileged user, seed its env, and enable+start it. It is read-only, so it is safe to
        # start immediately; it reads the daemon's local API (enabled by default in the seeded
        # config). It binds 0.0.0.0 by default — see dashboard.env to restrict it.
        if [ -e /usr/lib/easy-failover-dashboard/server.js ]; then
            if ! getent passwd easy-failover-dashboard >/dev/null 2>&1; then
                useradd --system --no-create-home --shell /usr/sbin/nologin \
                    --comment "easy-failover dashboard" easy-failover-dashboard || true
            fi
            if [ ! -e /etc/easy-failover-dashboard/dashboard.env ] \
                && [ -e /etc/easy-failover-dashboard/dashboard.env.example ]; then
                cp /etc/easy-failover-dashboard/dashboard.env.example \
                    /etc/easy-failover-dashboard/dashboard.env
                chmod 0644 /etc/easy-failover-dashboard/dashboard.env
            fi
        fi

        if command -v systemctl >/dev/null 2>&1; then
            # A stale SysV init script of the same name (never shipped by this package) makes
            # `systemctl enable` abort via update-rc.d ("Default-Start contains no runlevels").
            # If one exists and no dpkg package owns it, move it aside so the native unit enables.
            if [ -e /etc/init.d/easy-failover ] \
                && ! dpkg-query -S /etc/init.d/easy-failover >/dev/null 2>&1; then
                mv /etc/init.d/easy-failover \
                    "/etc/init.d/easy-failover.disabled-by-pkg.$(date +%s)" 2>/dev/null || true
            fi

            systemctl daemon-reload || true
            # Auto-start the daemon. It idles until a VIP is configured (clean-slate config),
            # so this is safe on a fresh install; configure the VIP/peers via the dashboard or
            # /etc/easy-failover/config.toml and it begins failover on the next restart.
            # Clear any failed/start-limited state from a previous version, then enable (for boot)
            # and start separately so an odd environment that trips `enable` still starts the unit.
            systemctl reset-failed easy-failover.service 2>/dev/null || true
            systemctl enable easy-failover.service 2>/dev/null || true
            systemctl start easy-failover.service || true
            if [ -e /usr/lib/easy-failover-dashboard/server.js ]; then
                systemctl reset-failed easy-failover-dashboard.service 2>/dev/null || true
                systemctl enable easy-failover-dashboard.service 2>/dev/null || true
                systemctl start easy-failover-dashboard.service || true
            fi
        fi
        ;;
esac

exit 0
