# packaging/deb/postinst
#!/bin/sh
# Debian/Ubuntu post-install for the easy-failover DAEMON package: seed config, provision a
# write token, start the service. The dashboard is a separate package with its own scripts.
set -e

TOKEN_FILE=/etc/easy-failover/api.token

case "$1" in
    configure)
        # On first install, create an active config from the shipped example if the
        # operator does not already have one. Upgrades never overwrite it.
        seeded_config=0
        if [ ! -e /etc/easy-failover/config.toml ] \
            && [ -e /etc/easy-failover/config.example.toml ]; then
            cp /etc/easy-failover/config.example.toml /etc/easy-failover/config.toml
            chmod 0644 /etc/easy-failover/config.toml
            seeded_config=1
        fi

        # Generate a random write token once (root-only). The daemon authenticates write
        # requests against it; the optional dashboard package reads the same file to wire its
        # Apply button without the operator plumbing a secret by hand.
        if [ ! -e "$TOKEN_FILE" ]; then
            ( umask 077
              if command -v openssl >/dev/null 2>&1; then
                  openssl rand -hex 32 > "$TOKEN_FILE"
              else
                  head -c 32 /dev/urandom | od -An -tx1 | tr -d ' \n' > "$TOKEN_FILE"
              fi )
            chmod 0600 "$TOKEN_FILE"
        fi
        TOKEN="$(cat "$TOKEN_FILE" 2>/dev/null || true)"

        # On a freshly seeded config, enable token-authenticated write mode (so a write client
        # like the dashboard can apply config). We never rewrite an operator's existing config.
        if [ "$seeded_config" = "1" ] && [ -n "$TOKEN" ]; then
            sed -i \
                -e 's|^[[:space:]]*read_only[[:space:]]*=.*|read_only = false|' \
                -e "s|^[[:space:]]*auth_token_file[[:space:]]*=.*|auth_token_file = \"${TOKEN_FILE}\"|" \
                /etc/easy-failover/config.toml || true
        fi

        if command -v systemctl >/dev/null 2>&1; then
            # A stale SysV init script of the same name (never shipped by this package) makes
            # `systemctl enable` abort via update-rc.d ("Default-Start contains no runlevels").
            # If one exists and no dpkg package owns it, move it aside so the native unit enables.
            if [ -e /etc/init.d/easy-failover ] \
                && ! dpkg-query -S /etc/init.d/easy-failover >/dev/null 2>&1; then
                mv /etc/init.d/easy-failover \
                    "/etc/init.d/easy-failover.disabled-by-pkg.$(date +%s)" 2>/dev/null || true
            fi

            systemctl daemon-reload || true
            # Auto-start the daemon. It idles until a VIP is configured (clean-slate config),
            # so this is safe on a fresh install. Clear any failed/start-limited state from a
            # previous version, then enable (for boot) and start separately so an odd environment
            # that trips `enable` still starts the unit.
            systemctl reset-failed easy-failover.service 2>/dev/null || true
            systemctl enable easy-failover.service 2>/dev/null || true
            systemctl start easy-failover.service || true
        fi
        ;;
esac

exit 0
